# roles/runelite/files/static-runelite-net.service
[Unit]
Description=Local static.runelite.net mirror

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/runefoil static
