# docker/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:tomcat]
command=/opt/tomcat8/bin/catalina.sh run
user=tomcat
environment=CATALINA_PID="/opt/tomcat8/temp/tomcat.pid",CATALINA_HOME="/opt/tomcat8",CATALINA_BASE="/opt/tomcat8"
autostart=false
autorestart=false
startsecs=2
stdout_logfile=/opt/tomcat8/logs/stdout.log
stderr_logfile=/opt/tomcat8/logs/stderr.log

[program:static_runelite_net]
command=/opt/runefoil/bin/runefoil-static-runelite-net
user=btw
autostart=false
autorestart=false
startsecs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
